# supervisor.conf
[program:web]
command=gunicorn3 -b localhost:8000 -w 4 frontend:web
directory=/home/ubuntu/ScrapeBot
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
